# etc/services.d/tinyproxy/run
#!/bin/sh
sleep 5  # wait a while. proxychains.conf should be filled up with proxies
proxychains tinyproxy -d -c /etc/tinyproxy.conf
