# etc/services.d/tinyproxy/run
#!/bin/sh
proxychains tinyproxy -d -c /etc/tinyproxy.conf
